# config/nvrr.service
[Unit]
Description=NVRR Backend
After=network.target mediamtx.service

[Service]
Type=simple
User=nvrr
WorkingDirectory=/opt/nvrr
ExecStart=/root/.local/bin/uv run python -m uvicorn main:app --host 127.0.0.1 --port 8000 --app-dir /opt/nvrr/backend
Restart=always
RestartSec=5

# Configuration
Environment=ADMIN_PASSWORD=changeme
Environment=NVRR_DB_PATH=/opt/nvrr/data/nvrr.db
Environment=MEDIAMTX_CONFIG_PATH=/opt/nvrr/mediamtx/mediamtx.yml
Environment=STREAM_MODE=rtsp
Environment=PATH=/root/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
